# profile.d/go.sh
#export GOPATH=/usr/share/go
#export GOPATH=$GOROOT/src:$GOROOT/pkg

#export GOROOT=/usr/lib/go-1.10

#export GOBIN=$GOROOT/bin:$HOME/go/bin:/usr/share/go/bin
#export PATH=$PATH:$GOBIN


export GOROOT=/usr/local/go
export GOPATH=$HOME/go

export PATH=$PATH:$GOROOT/bin:$GOPATH/bin
